# todo/.config/systemd/user/apagar.service
[Unit]
Description=Apagar PC con aviso

[Service]
Type=oneshot
ExecStart=/home/justin/scripts/apagar.sh
