# todo/.config/systemd/user/ejercicio.service
[Unit]
Description=Modo ejercicio

[Service]
Type=oneshot
ExecStart=/home/justin/.local/bin/ejercicio.sh
